# config/root/etc/systemd/system/automount-restarter@.service
[Unit]
Description=automount restarter for %i

[Service]
Type=oneshot
ExecStartPre=/usr/bin/sleep 10
ExecStart=/usr/bin/systemctl restart %i.automount
