# config/root/etc/systemd/system/mnt-Remote-NAS-Wiki.automount
[Unit]
Description=Automount NAS Wiki Share
OnFailure=automount-restarter@%N.service

# Network interface target dependancies `systemctl list-units --no-pager | grep net-device`
After=sys-subsystem-net-devices-WGHome.device
Requisite=sys-subsystem-net-devices-WGHome.device

[Automount]
Where=/mnt/Remote/NAS/Wiki

[Install]
WantedBy=multi-user.target
